# debian/etc/profile.d/hobot-environment.sh
export PATH="/usr/hobot/bin/:${PATH}"
export LD_LIBRARY_PATH="/usr/hobot/lib:/app/lib:/app/pub/lib:/middleware/lib:/middleware/pub/lib:/usr/hobot/lib/sensor:/system/lib:/system/usr/lib:/lib:${LD_LIBRARY_PATH}"

