# debian/DEBIAN/postinst
#!/bin/sh

ldconfig

systemctl enable cam-service.service
systemctl start cam-service.service

exit 0

