# debian/DEBIAN/prerm
#!/bin/sh

systemctl disable cam-service.service

exit 0
